# zeroconf/mychembl-upnp.conf
description "UPnP Broadcast daemon for myChEMBL" 
author "ndedman@ebi.ac.uk"

start on runlevel [2345] 
stop on runlevel [!2345] 

respawn 

exec /usr/bin/mychembl-upnp.py

