# configuration/beaker.ini
[program:beaker_gunicorn]
command=gunicorn -c gunicorn.conf beaker:app ; the program (relative uses PATH, can take args)
directory=/home/chembl/chembl_beaker/conf.d  ; directory to cwd to before exec (def no cwd)
priority=995                  ; the relative start priority (default 999)
autorestart=true        ; whether/when to restart (default: unexpected)
startsecs=60                   ; number of secs prog must stay running (def. 1)
stopsignal=TERM               ; signal used to kill process (default TERM)
stopwaitsecs=60               ; max num secs to wait b4 SIGKILL (default 10)
stopasgroup=true             ; send stop signal to the UNIX process group (default false)
killasgroup=true             ; SIGKILL the UNIX process group (def false)
user=chembl                   ; setuid to this UNIX account to run the program
