# configuration/ipynb.ini
[program:ipynb]
command=jupyter notebook
process_name=%(program_name)s
numprocs=1
numprocs_start=0
autostart=true
autorestart=true
startsecs=1
startretries=3
exitcodes=0,2
stopsignal=QUIT
stopwaitsecs=10
user=chembl
redirect_stderr=false
stdout_logfile=/tmp/ipynb.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stdout_capture_maxbytes=0
stdout_events_enabled=false
stderr_logfile=/tmp/ipynb_err.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
stderr_capture_maxbytes=0
stderr_events_enabled=false
environment=HOME="/home/chembl",SHELL="/bin/bash",USER="chembl",PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",PYTHONPATH="/opt/indigo"
directory=/home/chembl/ipynb_workbench
serverurl=AUTO
